# recipes-core/serial-console/serial-console/serial-console.conf
# serial-console - Start a getty on the serial line
#

start on local-filesystems

stop on runlevel [!2345]

respawn
exec /sbin/getty -L @SERIALPORT@ @BAUDRATE@ linux

