# recipes-core/images/core-image-apq/serial-console.conf
# serial-console - Start a getty on the serial line
#

start on runlevel [2345]

stop on runlevel [!2345]

respawn
exec /sbin/getty -L ttyHSL0 115200 linux

